# knot/rootfs/etc/services.d/babeltower/finish
#!/usr/bin/with-contenv bashio

$(bashio::log.info "Stoping babel tower")