# knot/rootfs/etc/services.d/babeltower/run
#!/usr/bin/with-contenv bashio

$(bashio::log.info "Starting babel tower")
RABBITMQ_IP=$(bashio::addon.ip_address "rabbitmq")
ENV=development RABBITMQ_URL=amqp://guest:guest@$RABBITMQ_IP USERS_HOSTNAME=localhost AUTH_HOSTNAME=localhost THINGS_HOSTNAME=localhost /usr/bin/babeltower
