# knot/rootfs/etc/services.d/postgres/run
#!/usr/bin/with-contenv bashio

su postgres -c 'pg_ctl start -D /var/lib/postgresql/data'
while ! psql -U "postgres" postgres -c "" 2> /dev/null; do
	sleep 1
done

sleep infinity